# apps/api/deploy/systemd/qassess-api.service
[Unit]
Description=QAssess API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qassess
Group=qassess
WorkingDirectory=/opt/qassess/apps/api
EnvironmentFile=/etc/qassess/api.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
